# chisel/rootfs/etc/s6-overlay/s6-rc.d/chisel/finish
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: Chisel
# Finishes the Chisel tunnel service
# ==============================================================================

bashio::log.info "Chisel tunnel service stopped" 