# chisel/rootfs/etc/s6-overlay/s6-rc.d/init-chisel-config/up
/bin/bash -c /etc/s6-overlay/scripts/chisel-config.sh 